# setup/gentoo2.6/nds/etc/init.d/nds_nds0
#!/sbin/runscript

opts="${opts}"
hostname=`hostname`

start() {
	ebegin "Starting NDS Server on nds0"
    # archive old log files, creating archive if needed
	cd /opt/rtcds/llo/l1/target/l1daqnds0
    mkdir -p logs/archive
	mv -f logs/nds.log*  logs/archive/

    # start up NDS server with date-stamped log
    stamp=`date +%s`
	start-stop-daemon --start --quiet -b -m --pidfile /var/run/nds.pid --chuid controls --exec /opt/rtcds/llo/l1/target/l1daqnds0/nds -- -l /opt/rtcds/llo/l1/target/l1daqnds0/logs/nds.log.${stamp} /opt/rtcds/llo/l1/target/l1daqnds0/pipe
	eend $?
}

stop() {
	ebegin "Stopping NDS Server on nds0"
	start-stop-daemon --stop --quiet --chuid controls --exec /opt/rtcds/llo/l1/target/l1daqnds0/nds
	eend $?
}
