# support/systemd/rts-awgtpman@.service
[Unit]
Description=Advanced LIGO RTS awgtpman: %i
After=rts-module@%i.service
Requisite=rts-module@%i.service
PartOf=rts@%i.target

[Service]
User=advligorts
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_SYS_NICE
AmbientCapabilities=CAP_SYS_NICE
EnvironmentFile=/etc/advligorts/systemd_env
ExecStart=/usr/bin/rts_awgtpman_exec %i
Restart=always
