# support/systemd/rts@.target
[Unit]
Description=Advanced LIGO RTS target: %i
Requires=rts-epics@%i.service rts-module@%i.service rts-awgtpman@%i.service 
