# setup/gentoo3.2/fe/etc/init.d/dis_nodemgr
#!/sbin/runscript
#  /etc/init.d/dis_nodemgr

# set node port based on table

DEFPORT=3443
PORTCFG=`/etc/dxportcfg.sh | awk '{print $2}'`
if [ ! $PORTCFG ]; then
    PORTCFG=$DEFPORT
fi

start() {
    ebegin "Starting Dolphin node manager on port $PORTCFG"
    start-stop-daemon --start --quiet --exec /opt/DIS/sbin/dis_nodemgr -- -basedir /opt/DIS -portnode $PORTCFG
    eend $?
}

stop() {
    ebegin "Stopping Dolphin node manager"
    start-stop-daemon --stop --quiet --exec /opt/DIS/sbin/dis_nodemgr -- -basedir /opt/DIS -portnode $PORTCFG
    eend $?
}

