# setup/gentoo3.2/fe/etc/local.d/local.stop
# remove DX node from list (if present)
if /etc/dolphin_dx_config.sh
then
    /opt/DIS/sbin/dxtool prepare-shutdown 0
    echo Dolphin DX prepared for shutdown
    sleep 3
fi
# remove IX node from list (if present)
if /etc/dolphin_ix_config.sh
then
    /opt/DIS/sbin/ixtool prepare-shutdown 0
    echo Dolphin IX prepared for shutdown
    sleep 3
fi
/etc/init.d/monit stop
pkill mx_stream
/opt/open-mx/sbin/omx_init stop
halt -p
