# setup/gentoo3.0/fe/etc/init.d/localmount
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need checkfs
}

start() {
	# Mount local filesystems in /etc/fstab.
	ebegin "Mounting local filesystems"
	mount -at noproc,noshm${NET_FS_LIST:+,no}${NET_FS_LIST// /,no} \
		-O no_netdev >/dev/null
	eend $? "Some local filesystem failed to mount"

	# Make sure we insert usbcore if its a module
	if [[ -f /proc/modules && ! -d /proc/bus/usb ]] ; then
		# >/dev/null to hide errors from non-USB users
		modprobe usbcore &> /dev/null
	fi
	
# Front-ends - disable USB FS

	# Setup Kernel Support for miscellaneous Binary Formats
	if [[ -d /proc/sys/fs/binfmt_misc ]] ; then
		local binfmt=$(grep -Fow binfmt_misc /proc/filesystems)
		if [[ -n ${binfmt} ]] ; then
			ebegin $"Mounting misc binary format filesystem"
			mount -t binfmt_misc binfmt_misc /proc/sys/fs/binfmt_misc \
				-o nodev,noexec,nosuid
			eend $?
		fi
	fi

	# We do our swapping here instead of rc so we can get urandom started
	# before us for people that like an encrypted swap.
	ebegin "Activating (possible) swap"
	/sbin/swapon -a
	eend $?

	# Start dm-crypt mappings, if any
	start_addon dm-crypt
}

# vim:ts=4
