# setup/gentoo3.0/fe/etc/init.d/mx_stream0
#!/sbin/runscript

opts="${opts}"
sys=`/etc/rt.sh`
hostname=`hostname`
if [ -x "${sys}" ] 
then
	exit 1
fi

# find line number in rtsystab. Use that to mx_stream slot on card (0-15)
#  Test stand has 4 FEs, so spread across cards
line_num=`grep -v ^# /etc/rtsystab | grep --perl-regexp -n "^${hostname}\s" | sed 's/^\([0-9]*\):.*/\1/g'`
line_off=$(expr $line_num - 1)
epcalc=$(expr $line_off % 2)
# add 1 to avoid end-point 0
epnum=$(expr $epcalc + 1)
cnum=$(expr $line_off / 2)

start() {
	ebegin "Starting mx_stream rt $line_num on ep $epnum card $cnum x2daqdc0"
	start-stop-daemon --start --quiet -b -m --pidfile /var/log/mx_stream0.pid --exec /opt/rtcds/tst/x2/target/x2daqdc0/mx_stream -- -e 0 -r "$epnum" -W 0 -w 0 -s "$sys" -d x2daqdc0:$cnum -l /opt/rtcds/tst/x2/target/x2daqdc0/mx_stream_logs/$hostname.log
	eend $?
}

stop() {
	ebegin "Stopping mx_stream on x2daqdc0"
	start-stop-daemon --stop --quiet --exec /opt/rtcds/tst/x2/target/x2daqdc0/mx_stream
	eend $?
}
