# setup/gentoo3.0/fe/etc/init.d/mx_stream1
#!/sbin/runscript

opts="${opts}"
sys=`/etc/rt.sh`
hostname=`hostname`
if [ -x "${sys}" ] 
then
	exit 1
fi

# find line number in rtsystab. Use that to mx_stream slot on card (0-15)
#  Test stand has up to 6 FEs, so spread across cards (3 on each card)
line_num=`grep -v ^# /etc/rtsystab | grep --perl-regexp -n "^${hostname}\s" | sed 's/^\([0-9]*\):.*/\1/g'`
line_off=$(expr $line_num - 1)
# end-point will run 0,1,2,0,1,2 
epcalc=$(expr $line_off % 3)
# add 1 to avoid end-point 0
epnum=$(expr $epcalc + 1)
# card num will run 0,0,0,1,1,1
cnum=$(expr $line_off / 3)

start() {
	ebegin "Starting mx_stream rt $line_num ep $epnum card $cnum to x2daqdc1"
	start-stop-daemon --start --quiet -b -m --pidfile /var/log/mx_stream1.pid --exec /opt/rtcds/tst/x2/target/x2daqdc1/mx_stream -- -e 1 -r "$epnum" -W 0 -w 0 -s "$sys" -d x2daqdc1:$cnum -l /opt/rtcds/tst/x2/target/x2daqdc1/mx_stream_logs/$hostname.log
	eend $?
}

stop() {
	ebegin "Stopping mx_stream to x2daqdc1"
	start-stop-daemon --stop --quiet --exec /opt/rtcds/tst/x2/target/x2daqdc1/mx_stream
	eend $?
}
