# support/udev/gpstime.rules
KERNEL=="gpstime", ACTION=="add", OWNER="advligorts", GROUP="advligorts", MODE="0664", RUN+="__SHARE__/calc_gps_offset.py -a"
