# support/udev/gpstime.rules
KERNEL=="gpstime", ACTION=="add", GROUP="controls", RUN+="__SHARE__/calc_gps_offset.py -a"
